# src/services/acix/indexserver/acix-index.service.in
[Unit]
Description=ARC cache index server
After=network.target local_fs.target

[Service]
Type=forking
PIDFile=/var/run/acix-index.pid
ExecStart=@pkgdatadir@/acix-index-start

[Install]
WantedBy=multi-user.target
