# src/services/acix/indexserver/arc-acix-index.service.in
[Unit]
Description=ARC cache index - index server
After=network.target local_fs.target

[Service]
Type=forking
PIDFile=/var/run/arc-acix-index.pid
ExecStart=/usr/share/arc/arc-acix-index-start

[Install]
WantedBy=multi-user.target
