# src/services/candypond/arc-candypond.service.in
[Unit]
Description=A-REX cache service
After=local_fs.target remote_fs.target

[Service]
Type=forking
PIDFile=/var/run/arched-cache-service.pid
ExecStart=@pkgdatadir@/arc-cache-service-start

[Install]
WantedBy=multi-user.target
