# src/services/gridftpd/arc-gridftpd.service.in
[Unit]
Description=ARC gridftpd
Documentation=man:gridftpd(8)
After=local_fs.target remote_fs.target

[Service]
Type=forking
PIDFile=/var/run/gridftpd.pid
ExecStart=/usr/share/arc/arc-gridftpd-start

[Install]
WantedBy=multi-user.target
