# src/services/a-rex/arc-arex.service.in
[Unit]
Description=ARC grid manager
After=local_fs.target remote_fs.target

[Service]
Type=forking
PIDFile=/var/run/arched-arex.pid
ExecStart=/usr/share/arc/arc-arex-start

[Install]
WantedBy=multi-user.target
