# src/services/a-rex/arc-arex.service.in
[Unit]
Description=ARC Remote EXecution service
After=local_fs.target remote_fs.target

[Service]
ExecStart=@prefix@/@pkgdatasubdir@/arc-arex-start -f
NotifyAccess=all

[Install]
WantedBy=multi-user.target
