# src/services/acix/scanner/arc-acix-scanner.service.in
[Unit]
Description=ARC cache index - scanner server
After=network.target local_fs.target

[Service]
Type=forking
PIDFile=/run/arc-acix-scanner.pid
ExecStart=/usr/share/arc/arc-acix-scanner-start

[Install]
WantedBy=multi-user.target
