# src/services/data-staging/arc-datadelivery-service.service.in
[Unit]
Description=A-REX datadelivery service
After=local_fs.target remote_fs.target

[Service]
Type=forking
PIDFile=/var/run/arched-datadelivery-service.pid
ExecStart=/usr/share/arc/arc-datadelivery-service-start

[Install]
WantedBy=multi-user.target
