# src/services/data-staging/arc-datadelivery-service.service.in
[Unit]
Description=A-REX datadelivery service
After=local_fs.target remote_fs.target

[Service]
ExecStart=@prefix@/@pkgdatasubdir@/arc-datadelivery-service-start -f
NotifyAccess=all

[Install]
WantedBy=multi-user.target
