# src/services/data-staging/arc-datadelivery-service.service.in
[Unit]
Description=A-REX datadelivery service
After=local-fs.target remote-fs.target

[Service]
ExecStart=@prefix@/@pkgdatasubdir@/arc-datadelivery-service-start
NotifyAccess=all
Type=forking

[Install]
WantedBy=multi-user.target
