# src/services/ldap-infosys/arc-infosys-ldap.service.in
[Unit]
Description=NorduGrid/ARC local resource information system
Requires=arc-infosys-ldap-slapd.service
After=arc-infosys-ldap-slapd.service
BindsTo=arc-infosys-ldap-slapd.service

[Service]
Type=forking
ExecStartPre=@prefix@/@pkgdatasubdir@/create-bdii-config
ExecStart=/bin/sh /var/run/arc/infosys/bdii-update.cmd
NotifyAccess=all

[Install]
WantedBy=multi-user.target
