# src/services/ldap-infosys/nordugrid-arc-ldap-infosys.in
#!/bin/bash
#
# Init file for the NorduGrid/ARC LDAP based Information System
#
# This script offers an alternative to the per-service init scripts
# It offers a convenient way to start all the infosys services in
# a signle command
#
# chkconfig: - 80 20
# description: NorduGrid/ARC Information system
#
######################################################################

### BEGIN INIT INFO
# Provides:          nordugrid-arc-ldap-infosys
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 2 3 4 5 6
# Short-Description: NorduGrid/ARC Information LDAP Information System
# Description:       Init file for the NorduGrid/ARC LDAP Information System
### END INIT INFO

prefix=${ARC_LOCATION:-@prefix@}

case "$1" in
    start | stop | restart | force-reload | reload | status | condrestart | try-restart)
	[ -x @initddir@/nordugrid-arc-slapd ] && \
	    @initddir@/nordugrid-arc-slapd "$1"
	[ -x @initddir@/nordugrid-arc-bdii ] && \
	    @initddir@/nordugrid-arc-bdii "$1"
	[ -x @initddir@/nordugrid-arc-egiis ] && \
	    @initddir@/nordugrid-arc-egiis "$1"
	exit 0
	;;
    *)
	echo "Usage: $0 {start|stop|restart|force-reload|reload|condrestart|try-restart|status}"
	exit 1
	;;
esac
