# src/services/ldap-infosys/nordugrid-arc-ldap-infosys.in
#!/bin/bash
#
# Init file for the NorduGrid/ARC LDAP based Information System
#
# This script offers an alternative to the per-service init scripts
# It offers a convenient way to start all the infosys services in
# a signle command
#
# chkconfig: - 80 20
# description: NorduGrid/ARC Information system
#
######################################################################

case "$1" in
    start | stop | restart | force-reload | reload | status | condrestart | try-restart)
	[ -x /etc/init.d/nordugrid-arc-slapd ] && \
	    /etc/init.d/nordugrid-arc-slapd "$1"
	[ -x /etc/init.d/nordugrid-arc-bdii ] && \
	    /etc/init.d/nordugrid-arc-bdii "$1"
	[ -x /etc/init.d/nordugrid-arc-egiis ] && \
	    /etc/init.d/nordugrid-arc-egiis "$1"
	;;
    *)
	echo "Usage: $0 {start|stop|restart|force-reload|reload|condrestart|try-restart|status}"
	exit 1
	;;
esac
