# src/hed/daemon/scripts/arched.service.in
[Unit]
Description=ARC HED service
After=local_fs.target remote_fs.target

[Service]
Type=forking
PIDFile=/var/run/arched.pid
ExecStart=/usr/share/arc/arched-start

[Install]
WantedBy=multi-user.target
