# src/services/acix/cacheserver/acix-cache.service.in
[Unit]
Description=ARC cache index cache server
After=network.target local_fs.target

[Service]
Type=forking
PIDFile=/var/run/acix-cache.pid
ExecStart=@pkgdatadir@/acix-cache-start

[Install]
WantedBy=multi-user.target
